# data/mercure/supervisor/mercure.conf
[program:mercure]
environment=JWT_KEY="jwtSecretMercurePixelForce123456TokenToken", CORS_ALLOWED_ORIGINS="https://pixelforce.fr", PUBLISH_ALLOWED_ORIGINS="https://pixelforce.fr", ADDR="pixelforce.fr:3000", MERCURE_EXTRA_DIRECTIVES="cors_origins https://pixelforce.fr", CORS_ORIGINS="https://pixelforce.fr"
command=/etc/mercure/mercure_v0_14_1/mercure
process_name=%(program_name)s_$(process_num)s
numprocs=1
autostart=true
#directory=/tmp
autorestart=true
startsecs=5
startretries=10
redirect_stderr=false
stdout_capture_maxbytes=1MB
stderr_capture_maxbytes=1MB
stdout_logfile=/etc/mercure/mercure_v0_14_1/mercureout.log
stderr_logfile=/etc/mercure/mercure_v0_14_1/mercureerror.log
